# priv/setup/services/chivi-bak.service
[Unit]
Description=chivi backup
After=network.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
User=nipin
WorkingDirectory=/app/chivi.app
ExecStart=/usr/bin/bash run/run-hourly.sh

[Install]
WantedBy=multi-user.target
