# priv/setup/services/cvapp-srv.service
[Unit]
Description=cvapp server
Requires=network.target postgresql.service
After=postgresql.service
AllowIsolate=yes

StartLimitIntervalSec=400
StartLimitBurst=20
StartLimitAction=reboot

[Service]
Type=simple
Restart=always
RestartSec=15
TimeoutSec=600

User=nipin
Environment=PORT=5010
Environment=CV_ENV=production
WorkingDirectory=/2tb/app.chivi
ExecStart=/2tb/app.chivi/bin/cvapp-srv

StandardOutput=file:/2tb/app.chivi/log/cvapp-srv.log
StandardError=append:/2tb/app.chivi/log/cvapp-err.log

[Install]
WantedBy=multi-user.target
