# priv/setup/services/cvmtl-mon.path
[Path]
Unit=cvmtl-mon.service
PathChanged=/app/chivi/var/cvmtl

[Install]
WantedBy=multi-user.target
