# priv/setup/services/cvmtl-mon.path
[Path]
Unit=mt_v2-mon.service
PathChanged=/app/chivi/var/mt_v2

[Install]
WantedBy=multi-user.target
