# priv/setup/services/mt_sp-dev.service
[Unit]
Description=mt_sp server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=nipin
Environment=PORT=5301
Environment=CV_ENV=development
WorkingDirectory=/2tb/dev.chivi
ExecStart=/2tb/dev.chivi/bin/mt_sp-srv
StandardOutput=file:/2tb/dev.chivi/log/mt_sp-srv.log
StandardError=append:/2tb/dev.chivi/log/mt_sp-err.log

[Install]
WantedBy=multi-user.target
