# priv/setup/services/mt_v2-mon.path
[Path]
Unit=mt_v2-mon.service
PathChanged=/app/chivi/var/cvmtl

[Install]
WantedBy=multi-user.target
